# scripts/chore-app.service
[Unit]
Description=Family Chore App
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/family-chore-app
ExecStart=/usr/bin/npm start
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3000

[Install]
WantedBy=multi-user.target
